# odoo/templates/supervisord.conf
[inet_http_server]
port=127.0.0.1:9001

[supervisord]
logfile=/opt/odoo/logs/supervisor/supervisord.log
logfile_maxbytes=20MB
logfile_backups=4
loglevel=debug
pidfile=/opt/odoo/logs/supervisor/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[program:odoo]
user=odoo
command=/opt/odoo/envs/odoo/bin/python /opt/odoo/odoo/odoo.py --addons-path=/opt/odoo/odoo/addons
directory=/opt/odoo/odoo/
